# legacy/systemd/system/drivespace.timer
[Unit]
Description=Run Capacity Report Service

[Timer]
OnBootSec=1min
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
