# legacy/systemd/system/lockfile.timer
[Unit]
Description=Run monitor_lock_file.service every 1 minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target