# storage-watchdog.timer
[Unit]
Description=Run storage-watchdog service every 30 minutes

[Timer]
OnBootSec=90
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
